# sv/irc/run
#!/bin/sh

exec urxvt -fn "xft:DejaVu Sans Mono:size=10" -b 8 -e weechat
