# src/sv/irc/run
#!/bin/sh

exec urxvt -fn "xft:DejaVu Sans Mono:size=16" -b 8 -e weechat
