# src/sv/irc/run
#!/bin/sh

exec urxvt -title "foubot9000:irc" -fn "xft:DejaVu Sans Mono:size=20" -b 8 -e weechat
