# src/sv/météo/run
#!/bin/sh

exec urxvt -fn "xft:DejaVu Sans Mono:size=22" -b 8 -e hwatch -ctx -n 300 -- run-wego.sh
